# configCopy/systemd/user/vdirsyncer.path
[Unit]
Description=Sync calendars when local vdirs change

[Path]
PathModified=%h/.calendars

[Install]
WantedBy=default.target
